# backend/deploy/securehub.service
[Unit]
Description=SecureHub FastAPI app
After=network.target

[Service]
User=securehub
Group=securehub
WorkingDirectory=/opt/securehub
Environment="PATH=/opt/securehub/venv/bin"
ExecStart=/opt/securehub/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=on-failure

[Install]
WantedBy=multi-user.target
